# ansible/roles/httpd/templates/etc/init/vagrant-httpd.conf.j2
start on vagrant-mounted

script
  /etc/init.d/httpd restart
end script
